# pi/kiosk/mood-plotter-kiosk.service
[Unit]
Description=Mood Plotter — Chromium-kiosk pa LCD:n
After=mood-plotter-wifi.service graphical.target
Wants=mood-plotter-wifi.service

[Service]
Type=simple
User=pi
PAMName=login
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
StandardInput=tty
StandardOutput=journal
StandardError=journal
UtmpIdentifier=tty1
UtmpMode=user
Environment=XDG_RUNTIME_DIR=/run/user/1000
# Lat wlroots ladda default-backenden (drm + libinput) — INTE bara drm.
# Den gamla "WLR_BACKENDS=drm" stangde av libinput, sa touch funkade aldrig.
Environment=WLR_NO_HARDWARE_CURSORS=1
Environment=XDG_SEAT=seat0
Environment=XDG_VTNR=1
# Cage = minimal Wayland-compositor som korar EN app fullskarm utan WM.
# Chromium i kiosk-lage pekar pa lokala Flask-appen.
ExecStart=/usr/bin/cage -s -- /usr/bin/chromium \
    --kiosk \
    --noerrdialogs \
    --disable-infobars \
    --disable-translate \
    --no-first-run \
    --start-fullscreen \
    --check-for-update-interval=31536000 \
    --overscroll-history-navigation=0 \
    --disable-pinch \
    --disable-features=TranslateUI \
    --hide-scrollbars \
    --autoplay-policy=no-user-gesture-required \
    --ozone-platform=wayland \
    --enable-features=UseOzonePlatform \
    --touch-events=enabled \
    --window-size=800,480 \
    http://127.0.0.1:5000/
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
