# pi/kiosk/mood-plotter-kiosk.service
[Unit]
Description=Mood Plotter — Sway-kiosk pa LCD:n
# udev-settle: vanta tills 99-waveshare-touch.rules har applicerats sa
# FT5x06 ar klassad som touchscreen INNAN sway startar. Utan detta kan
# sway plocka upp enheten som mus i forsta forsoket → ingen touch.
After=mood-plotter-wifi.service graphical.target systemd-udev-settle.service
Wants=mood-plotter-wifi.service systemd-udev-settle.service

[Service]
Type=simple
User=pi
PAMName=login
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
StandardInput=tty
StandardOutput=journal
StandardError=journal
UtmpIdentifier=tty1
UtmpMode=user
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=WLR_NO_HARDWARE_CURSORS=1
Environment=XDG_SEAT=seat0
Environment=XDG_VTNR=1
# Sway istallet for cage: sway har 'seat * hide_cursor 0' som
# definitivt doljer muspekaren. Cage saknar den featuren och dess
# default-cursor-fallback gick inte att doverstyra via XCURSOR_THEME.
ExecStart=/usr/bin/sway -c /home/pi/mood-plotter/pi/kiosk/sway-kiosk.conf
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
